# etc/xinetd.d/busybox-telnetd
service telnet
{
	socket_type = stream
	flags = IPv6
	user = root
	server = /usr/sbin/telnetd
	wait = no
}
